# projects/dashboard/systemd-dashboard.service
[Unit]
Description=系统监控看板服务
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/usr/share/nginx/html/dashboard
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/usr/share/nginx/html/dashboard/backend"
Environment="FLASK_ENV=production"
Environment="FLASK_PORT=5000"

# 启动命令
ExecStart=/usr/share/nginx/html/dashboard/start.sh start
ExecStop=/usr/share/nginx/html/dashboard/start.sh stop
ExecReload=/usr/share/nginx/html/dashboard/start.sh restart

# 重启策略
Restart=always
RestartSec=5

# 日志输出
StandardOutput=append:/var/log/dashboard/service.log
StandardError=append:/var/log/dashboard/error.log

[Install]
WantedBy=multi-user.target
